# config/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper every 8 hours

[Timer]
OnActiveSec=0sec
OnUnitActiveSec=8h
AccuracySec=500msec

[Install]
WantedBy=X.target
