# config/.config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper at regular intervals

[Service]
Type=oneshot
Environment=DISPLAY=:0.0
ExecStart=/home/geertje/.bin/gwall
